# queue.conf
[program:queue]
command=/usr/local/bin/php /var/www/artisan queue:work
numprocs=1
autostart=true
autorestart=true
user=root